# debian/onenas-installer.service
[Unit]
Description=TrueNAS Installer

[Service]
ExecStart=exec fbterm -- /usr/bin/python3 -m onenas_installer

[Install]
WantedBy=multi-user.target
